# config/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = /usr/bin/uwsgi --ini /opt/webhook/config/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off; user nginx nginx;'
